# meta-bsp-rpi/recipes-core/udev/udev-rules-rpi/vcan.rules
SUBSYSTEM=="module", KERNEL=="vcan", ACTION=="add", RUN+="/etc/udev/scripts/vcan.sh add"
SUBSYSTEM=="module", KERNEL=="vcan", ACTION=="remove", RUN+="/etc/udev/scripts/vcan.sh remove"